# system/lyrah-first-boot.service
[Unit]
Description=Lyrah OS First Boot Setup
ConditionPathExists=!/var/lib/lyrah/.first-boot-done

[Service]
Type=oneshot
ExecStart=/usr/share/lyrah/setup/first-boot.sh
RemainAfterExit=yes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
